# supervisor/watermon.conf
[program:watermon]
command=/home/pi/watermon/start.sh
numprocs=1
directory=/home/pi/watermon
autostart=true
user=pi
